# filesystems/vnx_rootfs_lxc_ubuntu64-16.04-v025-openstack-compute/rootfs/var/lib/dpkg/info/libvirt-daemon-system.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/libvirt-daemon-system
/etc/default
/etc/default/virtlogd
/etc/default/libvirt-guests
/etc/default/libvirtd
/etc/logrotate.d
/etc/logrotate.d/libvirtd.uml
/etc/logrotate.d/libvirtd.qemu
/etc/logrotate.d/libvirtd.lxc
/etc/logrotate.d/libvirtd
/etc/logrotate.d/libvirtd.libxl
/etc/dnsmasq.d-available
/etc/dnsmasq.d-available/libvirt-daemon
/etc/sasl2
/etc/sasl2/libvirt.conf
/etc/apparmor.d
/etc/apparmor.d/usr.sbin.libvirtd
/etc/apparmor.d/abstractions
/etc/apparmor.d/abstractions/libvirt-qemu
/etc/apparmor.d/abstractions/libvirt-lxc
/etc/apparmor.d/local
/etc/apparmor.d/local/usr.sbin.libvirtd
/etc/apparmor.d/local/usr.lib.libvirt.virt-aa-helper
/etc/apparmor.d/usr.lib.libvirt.virt-aa-helper
/etc/apparmor.d/libvirt
/etc/apparmor.d/libvirt/TEMPLATE.lxc
/etc/apparmor.d/libvirt/TEMPLATE.qemu
/etc/init.d
/etc/init.d/virtlogd
/etc/init.d/libvirt-guests
/etc/init.d/libvirtd
/etc/init.d/virtlockd
/etc/libvirt
/etc/libvirt/libxl.conf
/etc/libvirt/virtlogd.conf
/etc/libvirt/qemu-sanlock.conf
/etc/libvirt/lxc.conf
/etc/libvirt/qemu
/etc/libvirt/qemu/networks
/etc/libvirt/qemu/networks/default.xml
/etc/libvirt/qemu/networks/autostart
/etc/libvirt/hooks
/etc/libvirt/nwfilter
/etc/libvirt/nwfilter/allow-arp.xml
/etc/libvirt/nwfilter/no-arp-mac-spoofing.xml
/etc/libvirt/nwfilter/qemu-announce-self-rarp.xml
/etc/libvirt/nwfilter/no-arp-ip-spoofing.xml
/etc/libvirt/nwfilter/no-other-rarp-traffic.xml
/etc/libvirt/nwfilter/allow-ipv4.xml
/etc/libvirt/nwfilter/no-mac-broadcast.xml
/etc/libvirt/nwfilter/allow-dhcp-server.xml
/etc/libvirt/nwfilter/allow-dhcp.xml
/etc/libvirt/nwfilter/no-other-l2-traffic.xml
/etc/libvirt/nwfilter/no-mac-spoofing.xml
/etc/libvirt/nwfilter/no-ip-spoofing.xml
/etc/libvirt/nwfilter/no-ip-multicast.xml
/etc/libvirt/nwfilter/clean-traffic.xml
/etc/libvirt/nwfilter/no-arp-spoofing.xml
/etc/libvirt/nwfilter/allow-incoming-ipv4.xml
/etc/libvirt/nwfilter/qemu-announce-self.xml
/etc/libvirt/libvirtd.conf
/etc/libvirt/virtlockd.conf
/etc/libvirt/libxl-lockd.conf
/etc/libvirt/qemu-lockd.conf
/etc/libvirt/qemu.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/virt-guest-shutdown.target
/lib/systemd/system/libvirt-guests.service
/lib/systemd/system/virtlogd.service
/lib/systemd/system/virtlockd.socket
/lib/systemd/system/virtlogd.socket
/lib/systemd/system/libvirtd.service
/lib/systemd/system/virtlockd.service
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/60-libvirt.pkla
/var/lib/libvirt
/var/lib/libvirt/images
/var/lib/libvirt/qemu
/var/lib/libvirt/qemu/channel
/var/lib/libvirt/qemu/channel/target
/var/lib/libvirt/boot
/var/lib/libvirt/sanlock
/var/cache
/var/cache/libvirt
/var/cache/libvirt/qemu
/var/log
/var/log/libvirt
/var/log/libvirt/lxc
/var/log/libvirt/libxl
/var/log/libvirt/qemu
/var/log/libvirt/uml
/usr
/usr/lib
/usr/lib/libvirt
/usr/lib/libvirt/virt-aa-helper
/usr/share
/usr/share/bug
/usr/share/bug/libvirt-daemon-system
/usr/share/bug/libvirt-daemon-system/presubj
/usr/share/doc
/usr/share/doc/libvirt-daemon-system
/usr/share/doc/libvirt-daemon-system/examples
/usr/share/doc/libvirt-daemon-system/examples/libvirtd.conf.gz
/usr/share/doc/libvirt-daemon-system/examples/libvirtd.policy
/usr/share/doc/libvirt-daemon-system/examples/libvirtd.sasl
/usr/share/doc/libvirt-daemon-system/copyright
/usr/share/doc/libvirt-daemon-system/NEWS.Debian.gz
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/polkit-1
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/60-libvirt.rules
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.libvirt.unix.policy
/usr/share/polkit-1/actions/org.libvirt.api.policy
/usr/share/systemtap
/usr/share/systemtap/tapset
/usr/share/systemtap/tapset/libvirt_functions.stp
/usr/share/systemtap/tapset/libvirt_qemu_probes.stp
/usr/share/systemtap/tapset/libvirt_probes.stp
/usr/share/doc/libvirt-daemon-system/changelog.Debian.gz
/usr/share/doc/libvirt-daemon-system/README.Debian.gz
/usr/share/doc/libvirt-daemon-system/changelog.gz
