# filesystems/vnx_rootfs_lxc_ubuntu64-16.04-v025-openstack-compute/rootfs/var/lib/dpkg/info/neutron-openvswitch-agent.conffiles
/etc/init.d/neutron-openvswitch-agent
/etc/init.d/neutron-ovs-cleanup
/etc/init/neutron-openvswitch-agent.conf
/etc/init/neutron-ovs-cleanup.conf
/etc/neutron/plugins/ml2/openvswitch_agent.ini
/etc/neutron/rootwrap.d/openvswitch-plugin.filters
