# filesystems/vnx_rootfs_lxc_ubuntu64-16.04-v025-openstack-compute/rootfs/var/lib/dpkg/info/neutron-openvswitch-agent.list
/.
/usr
/usr/bin
/usr/bin/neutron-openvswitch-agent
/usr/share
/usr/share/doc
/usr/share/doc/neutron-openvswitch-agent
/usr/share/doc/neutron-openvswitch-agent/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/neutron-openvswitch-agent.8.gz
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/neutron-ovs-cleanup.service
/lib/systemd/system/neutron-openvswitch-agent.service
/etc
/etc/init
/etc/init/neutron-ovs-cleanup.conf
/etc/init/neutron-openvswitch-agent.conf
/etc/neutron
/etc/neutron/rootwrap.d
/etc/neutron/rootwrap.d/openvswitch-plugin.filters
/etc/neutron/plugins
/etc/neutron/plugins/ml2
/etc/neutron/plugins/ml2/openvswitch_agent.ini
/etc/init.d
/etc/init.d/neutron-ovs-cleanup
/etc/init.d/neutron-openvswitch-agent
/usr/share/doc/neutron-openvswitch-agent/changelog.Debian.gz
