# filesystems/vnx_rootfs_lxc_ubuntu64-16.04-v025-openstack-compute/rootfs/var/lib/dpkg/info/nova-compute.conffiles
/etc/init.d/nova-compute
/etc/init/nova-compute.conf
/etc/nova/rootwrap.d/compute.filters
