# filesystems/vnx_rootfs_lxc_ubuntu64-16.04-v025-openstack-compute/rootfs/var/lib/dpkg/info/nova-compute.list
/.
/usr
/usr/bin
/usr/bin/nova-compute
/usr/bin/nova-idmapshift
/usr/share
/usr/share/doc
/usr/share/doc/nova-compute
/usr/share/doc/nova-compute/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/nova-compute.1.gz
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nova-compute.service
/etc
/etc/init
/etc/init/nova-compute.conf
/etc/init.d
/etc/init.d/nova-compute
/etc/nova
/etc/nova/rootwrap.d
/etc/nova/rootwrap.d/compute.filters
/usr/share/doc/nova-compute/changelog.Debian.gz
