# filesystems/vnx_rootfs_lxc_ubuntu64-16.04-v025-openstack-compute/rootfs/var/lib/dpkg/info/openvswitch-switch.conffiles
/etc/default/openvswitch-switch
/etc/init.d/openvswitch-switch
/etc/init/openvswitch-force-reload-kmod.conf
/etc/init/openvswitch-switch.conf
/etc/logrotate.d/openvswitch-switch
